# controls/mj60pi/conf.d/dragonfly.conf
[program:mj60pi]
command = /home/pi/cage_venv/bin/dragonfly serve -vvv -c gpio.yaml
directory = /home/pi/cage/controls/mj60pi
startsecs = 10
stopsignal = INT
stdout_logfile_maxbytes = 5MB

[program:mj60ln]
command = /home/pi/cage_venv/bin/dragonfly serve -vvv -c lnscale.yaml
directory = /home/pi/cage/controls/mj60pi
startsecs = 10
stopsignal = INT
stdout_logfile_maxbytes = 5MB

