# controls/testpi/conf.d/dragonfly.conf
[program:testpi]
command = /home/pi/cage_venv/bin/dragonfly serve -vvv -c heaterpi.yaml
directory = /home/pi/cage/controls/testpi
startsecs = 10
stopsignal = INT
stdout_logfile_maxbytes = 5MB
